# ha-llm-tools/rootfs/etc/services.d/llm/run
#!/usr/bin/with-contenv bash
# ==============================================================================
# Start the HA LLM Tools service
# ==============================================================================

exec /usr/bin/start.sh